# rpmconf/Env/crontabs/crontab.logger
#
# crontab.logger
#
# process logs
#
00,10,20,30,40,50 * * * * /opt/zextras/libexec/zmlogprocess > /tmp/logprocess.out 2>&1
